# debian/postrm
#!/bin/bash
set -e

systemctl disable rp-server.service || true
systemctl stop rp-server.service || true
